# configs/rc.local
#!/bin/sh -e
cd /home/pi/HiveMind_Node
python HiveMind_Node.py configs/node.conf
exit 0
